# lxc/archlinux/myservice.service
[Unit]
Description=myservice
After=systemd-journald.service

[Service]
ExecStart=/myservice

[Install]
WantedBy=multi-user.target

