# docker/noble_novnc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:sssd]
command=/usr/sbin/sssd -i --logger=files
autostart=true
autorestart=true

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1024x768x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -rfbauth /root/.vnc/passwd -forever -shared -rfbport 5900
#command=/usr/bin/x11vnc -display :0 -auth guess -pam_service vnc -forever -shared -rfbport 5900
#command=/bin/true
autorestart=false

#[program:tigervnc]
#command=/usr/bin/x0vncserver -display :0 -rfbport 5900 -SecurityTypes=TLSVnc -PAMService=login
#autorestart=true
#priority=20

#[program:novnc]
#command=/usr/share/novnc/utils/novnc_proxy \
#  --vnc localhost:5900 \
#  --listen 8443 \
#  --cert /etc/ssl/certs/server.crt \
#  --key  /etc/ssl/private/server.key
#autorestart=true
#priority=30

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc --cert /etc/ssl/certs/server.crt --key  /etc/ssl/private/server.key 8443 localhost:5900
autorestart=true
priority=30

[program:xfce4]
command=/bin/bash -c "export DISPLAY=:0 && startxfce4"
autorestart=true
priority=40
environment=DISPLAY=":0"

