# docker/noble/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:sssd]
command=/usr/sbin/sssd -i --logger=files
autostart=true
autorestart=true

